# roles/nginx/files/certbot-renew.timer
[Unit]
Description=certbot-renew timer

[Timer]
Unit=certbot-renew.service
OnCalendar=daily

[Install]
WantedBy=timers.target
