# roles/nginx/files/certbot-renew.timer
[Unit]
Description=Monthly letsencrypt certificate renewal

[Timer]
Unit=certbot-renew.service
OnCalendar=monthly

[Install]
WantedBy=timers.target
