# roles/auto-update/files/apt-update.timer
[Unit]
Description=Monthly apt package update

[Timer]
Unit=apt-update.service
OnCalendar=*-*-01 05:00:00

[Install]
WantedBy=timers.target
