# roles/aws-ddns/templates/aws-ddns.timer.j2
[Unit]
Description=Daily update AWS dynamic DNS record

[Timer]
Unit=aws-ddns.service
OnCalendar=daily

[Install]
WantedBy=timers.target
