# roles/aws-ddns/templates/aws-ddns.timer.j2
[Unit]
Description=AWS DDNS timer

[Timer]
Unit=aws-ddns.service
OnCalendar=daily

[Install]
WantedBy=timers.target
